# rpm/log_memory_current.timer
[Unit]
Description=WMP periodic log of memory consumption

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
