# rpm/sap.slice
[Unit]
Description=WMP slice
Documentation=Resource controls are applied primarily to this slice.

