# rpm/sap.slice
# Resource controls are applied primarily to this slice.
[Unit]
Description=WMP slice

