# files/etc/systemd/system/sapinit.service.d/10-increase-limits.conf
[Service]
# To prevent killing of long-running termanating processes (e.g. databases) at shutdown.
TimeoutStopSec=30m

# Increasing TasksMax for big (HANA) systems.
TasksMax=infinity


