# files/etc/systemd/system/sapinit.service.d/10-sap.slice.conf
[Service]
# This file makes sure, that the `sapinit.service` slice will be moved from
# `system.slice` beneath `sap.slice`. Therefore all processes are covered by
# the memory protection too, which normally remain in the slice
# `sapinit.service`.
Slice=sap.slice

# To prevent killing of long-running termanating processes (e.g. databases) at shutdown.
TimeoutStopSec=30m

# Increasing TasksMax for big (HANA) systems.
TasksMax=infinity


