# files/etc/systemd/system/sapinit.service.d/50-MemoryLow.conf
[Service]
# Same value as in the parent slice (could be "max" and rely on parent's restriction)
MemoryLow=90%
